# modules/redis/templates/redis_server_init_script.conf.erb
# redis-server - Redis Datastore Server
#
# Redis is a key value in memory persistent datastore

start on (local-filesystems and runlevel [2345])
stop on runlevel [016]
respawn

limit nofile 20000 65000
pre-start script
    mkdir -p /var/run/redis
    chown redis:redis /var/run/redis
end script
<% if has_variable?("socket") and @socket != '' then %>
pre-stop script
    /usr/bin/redis -s <%= @socket %> SHUTDOWN
end script
<% end %>
<% if has_variable?("port") and @port != '' then %>
pre-stop script
    /usr/bin/redis -p <%= @port %> SHUTDOWN
end script
<% end %>
# with older upstart versions, this is the recomended way to change the user:group
exec start-stop-daemon --start --chuid redis:redis --pidfile <%= @pidfile %> --umask 007 --exec /usr/bin/redis-server -- <%= @conffile %>

